# data/buxton.socket.in
[Unit]
Description=Buxton Configuration Service

[Socket]
ListenStream=@BUXTON_SOCKET@

[Install]
WantedBy=sockets.target
