# data/buxton.service.in
[Unit]
Description=Buxton Configuration Service

[Service]
ExecStart=@prefix@/sbin/buxtond
User=@BUXTON_USERNAME@

[Install]
WantedBy=multi-user.target
