# data/buxton.socket.in
[Unit]
Description=Buxton Configuration Service

[Socket]
ListenStream=/run/buxton-0

[Install]
WantedBy=sockets.target
