# data/buxton.service.in
[Unit]
Description=Buxton Configuration Service

[Service]
ExecStart=/usr/bin/bt-daemon

[Install]
WantedBy=multi-user.target
